# tools/irc/cia.crontab
@reboot cd /home/commits/tools/irc; nohup ./announceBot.py&
0 0 * * * rm /home/commits/stats/end-of-day/*; mv /home/commits/stats/daily/* /home/commits/stats/end-of-day/
0 0 1 * * rm /home/commits/stats/monthly/*
0 0 * * 0 rm /home/commits/stats/weekly/*


